# config/launchd/com.journal-scout.ft50.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.journal-scout.ft50</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{INSTALL_PATH}/pipelines/ft50-daily.sh</string>
    </array>

    <!-- Daily at 9:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>{LOG_PATH}/launchd-ft50-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{LOG_PATH}/launchd-ft50-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{HOME}</string>
        <key>PATH</key>
        <string>{HOME}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
